# deploy/supervisor/erp-acesso-worker.conf
; Supervisor — fila Laravel (QR Code, e-mails, jobs)
; Instalado por: deploy/scripts/instalar-servicos.sh

[program:erp-acesso-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ERP-Acesso/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/ERP-Acesso/storage/logs/worker.log
stopwaitsecs=3600
