# deploy/supervisor/laravel-worker.conf
; Supervisor — fila Laravel (QR Code, e-mails, jobs futuros)
; Copiar para: /etc/supervisor/conf.d/linha-leve-worker.conf
;
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status linha-leve-worker:*

[program:linha-leve-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/linha-leve/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/linha-leve/storage/logs/worker.log
stopwaitsecs=3600
